# packer/redis-3.2.8/etc/systemd/navoptenv.path
[Unit]
Description=Wait for /etc/navoptenv.json to exist

[Path]
PathExists=/etc/navoptenv.json

[Install]
WantedBy=multi-user.target
