# packer/redis-3.2.8/etc/systemd/redis-backup.timer
[Unit]
Description=Redis backup
After=redis-server.service

[Timer]
OnCalendar=*:0/30

[Install]
WantedBy=multi-user.target
