# packer/elasticsearch/etc/cron.d/elasticsearch
# Set up hourly backup crontab
*/30 * * * * root /bin/bash /usr/local/bin/take_snapshot.sh >>/var/log/take_snapshot.log 2>&1
