# playbooks/roles/salt/templates/cron.d/salt-pull-config-apply
#Minute  Hour  DayOfMonth  Month    DayOfWeek  Command
1,16,31,46 * * * * root /sre/bin/salt_cron.sh
