# scripts/build-navopt/navopt-javaservices/advanalytics.conf
[supervisord]
nodaemon=true

[program:compileserver]
command=/usr/local/bin/xplaincompileserver
autorestart=true

[program:pyservice]
command=python /usr/local/bin/XplainAdvAnalyticsService.py